# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/nginx/stdout.log
stderr_logfile=/var/log/nginx/stderr.log

[program:flask]
command=/usr/bin/python3 /app/app.py
directory=/app
autostart=true
autorestart=true
environment=PYTHONUNBUFFERED="1"
stdout_logfile=/var/log/flask/stdout.log
stderr_logfile=/var/log/flask/stderr.log
